# dot_config/containers/systemd/caddy.container
# ~/.config/containers/systemd/caddy.container
# Quadlet unit — converts to a systemd service automatically

[Container]
Image=docker.io/library/caddy:latest
ContainerName=caddy
Environment=HTTP_PROXY=
Environment=HTTPS_PROXY=
Environment=NO_PROXY=
Environment=http_proxy=
Environment=https_proxy=
Environment=no_proxy=

PublishPort=443:443
PublishPort=80:80

# 挂载 Caddyfile 和数据卷
Volume=%h/.config/containers/caddy/Caddyfile:/etc/caddy/Caddyfile:ro
Volume=caddy-data:/data
Volume=caddy-config:/config

# 加入 Caddy 网络
Network=caddy-net
SecurityLabelDisable=true

[Unit]
After=network.target
[Service]
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=default.target
