# dot_config/containers/systemd/hindsight.container
# ~/.config/containers/systemd/hindsight.container
# Quadlet unit — converts to a systemd service automatically
# Docs: https://docs.podman.io/en/latest/markdown/podman-systemd.unit.5.html

[Container]
Image=ghcr.io/vectorize-io/hindsight:latest
ContainerName=hindsight
EnvironmentFile=%h/.config/containers/hindsight.env
UserNS=keep-id

# REST API + MCP（保留端口暴露便于调试）
PublishPort=8888:8888
# Web UI
PublishPort=9999:9999

# 加入 Caddy 网络
Network=caddy-net

# 数据目录 — 用 bind mount，方便备份和 rsync
Volume=%h/Documents/AppData/hindsight/data:/home/hindsight/.pg0

# 容器内进程以非 root 运行时的安全标签
SecurityLabelDisable=true

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
