# dot_config/containers/systemd/hindsight.container
# ~/.config/containers/systemd/hindsight.container
# Quadlet unit — converts to a systemd service automatically
# Docs: https://docs.podman.io/en/latest/markdown/podman-systemd.unit.5.html

[Container]
Image=ghcr.io/vectorize-io/hindsight:latest
ContainerName=hindsight
Environment=HTTP_PROXY=
Environment=HTTPS_PROXY=
Environment=NO_PROXY=
Environment=http_proxy=
Environment=https_proxy=
Environment=no_proxy=
EnvironmentFile=%h/.config/containers/hindsight.env
UserNS=keep-id

# REST API + MCP
PublishPort=8888:8888
# Web UI
PublishPort=9999:9999

Network=caddy-net

Volume=%h/Documents/AppData/hindsight/data:/home/hindsight/.pg0

# 容器内进程以非 root 运行时的安全标签
SecurityLabelDisable=true

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
